# etc/systemd/system/reload-php-fpm.timer
[Unit]
Description=Run reload PHP 7.4 FPM every day at 12 AM

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target